# Task_Notifier/Screen_Time/config/screen_time_script.timer
[Unit]
Description=Screen Time Notifier Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=timeperexecute

[Install]
WantedBy=wantedby