# Task_Notifier/Screen_Time/config/screen_time_script.service
[Unit]
Description=Screen Time Notifier

[Service]
Type=simple
ExecStart=execstart
Restart=on-failure

[Install]
WantedBy=wantedby